# sysadmin/server-etc-systemd-acme.sh.timer
[Unit]
Description=Daily check for acme.sh

[Timer]
OnCalendar=daily
RandomizedDelaySec=1h

[Install]
WantedBy=timers.target
