# sysadmin/server-opt-backup-all.service
[Unit]
Description=Backup all

[Service]
ExecStart=/opt/scripts/bin/backup-all
