# sway-session.target
[Unit]
Description=Sway Session
Before=graphical-session.target
BindsTo=graphical-session.target

BindsTo=sway.service
