# systemd-fetch-repos.service
[Unit]
Description=Fetch my repos

[Service]
ExecStart=%h/olv-etc/scripts/my-fetch-repos
