# admin/server-cron-backup-all.timer
[Unit]
Description=Backup everything weekly

[Timer]
OnStartupSec=1w
OnUnitActiveSec=1w

[Install]
WantedBy=timers.target
