# admin/server-opt-duckdns.service
[Unit]
Description=Update Duck DNS records

[Service]
ExecStart=/opt/scripts/bin/duckdns
User=notroot
