# admin/server-cron-duckdns.timer
[Unit]
Description=Update Duck DNS records periodically

[Timer]
OnStartupSec=5m
OnUnitActiveSec=1h

[Install]
WantedBy=timers.target
